# service-with-cron/crontab
#run python script every minutes
*/1 * * * * /usr/local/bin/python /service/cron_functions/get_smth.py > /proc/1/fd/1 2>/proc/1/fd/2
